# analysisd.service.template
[Unit]
Description=Protocol Analysis Daemon

[Service]
ExecStart=<execute-start> <config-file>
Type=simple
KillMode=control-group
Restart=on-failure
RestartSec=10s
